# supervisor/conf.d/aria2.conf
[program:aria2]
command=/usr/local/bin/aria2c --conf-path=/Users/samliu/dotfiles/aria2/aria2.conf

